# opnhost.service
[Unit]
Description=OpnHost

[Service]
User=root
WorkingDirectory=/OpnHost
User=root
ExecStart=/OpnHost/opnhosts

[Install]
WantedBy=multi-user.target